# osquery.conf
{
  "options": {
    "logger_snapshot_event_type": "true",
    "host_identifier": "hostname",
    "audit_allow_sockets": "true",
    "audit_allow_config": "true",
    "disable_audit": "false",
     "enable_syslog":"true",
     "disable_events":"false",
     "enable_bpf_events":"true",
      "enable_file_events":"true",
      "worker_threads": "4",
      "schedule_splay_percent": 10


  },
  "platform": "linux",
  "schedule": {

         "hash_processes": {
      "query" : "select p.pid , p.name , p.path , p.cmdline , p.state , p.cwd , p.root , p.parent, h.path, h.md5, h.sha1, h.sha256 , u.username  from  processes p join users u   on p.uid=p.uid join hash as  h on p.path=h.path  ;",
      "interval" : "84600",
      "platform": "posix",
      "version" : "0.0.1",
      "description" : "hash and md5  of processe ."
    },


         "parentname": {
      "query" : "select  p.pid , p.name, p.path, p.cmdline , p.cwd , p.root , p.parent, p.state, username, pp.name as parentname, p.name as processname , p.cmdline as processcmdline,  pp.cmdline as parentcmdline from  processes p join users on p.uid=users.uid join processes pp on  p.parent = pp.pid  ;",
      "interval" : "84600",
      "platform": "posix",
      "version" : "0.0.1",
      "description" : "parent name process name ."
    },


       "@crontab": {
      "query": "SELECT * FROM crontab;",
      "interval": 3600,
      "description": "Retrieves all the jobs scheduled in crontab in the target system."
    },

        "@dns_resolvers": {
      "query": "SELECT * FROM dns_resolvers;",
      "interval": 3600,
      "description": "DNS resolvers used by the host"
    },

 "@etc_hosts": {
      "query": "SELECT * FROM etc_hosts;",
      "interval": 3600,
      "description": "Retrieves all the entries in the target system /etc/hosts file."
    },

     "@iptables": {
      "query": "SELECT * FROM iptables;",
      "interval": 28800,
      "platform": "linux",
      "description": "Retrieves the current filters and chains per filter in the target system."
    },

    "@last": {
      "query": "SELECT * FROM last;",
      "interval": 84600,
      "description": "Retrieves the list of the latest logins with PID, username and timestamp."
    },

    "@user_ssh_keys": {
      "query": "SELECT * FROM users CROSS JOIN user_ssh_keys USING (uid);",
      "interval": 3600,
      "description": "Returns the private keys in the users ~/.ssh directory and whether or not they are encrypted"
    },


    "file_paths": {
    "home": [ "sig_group_1" ],
    "etc": [ "sig_group_1" ],
    "www": [ "sig_group_1" ],
    "binaries": [ "sig_group_1" ]

    }

    },


  "packs": {
       	"incident-response": "/opt/osquery/share/osquery/packs/incident-response.conf",
	    "linux_command_execution": "/opt/osquery/share/osquery/packs/linux_command_execution.conf",
        "fim": "/opt/osquery/share/osquery/packs/fim.conf"




  }
}

